# crawler.timer
[Unit]
Description=Run unified security crawler periodically

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
